# contestscore.service
[Unit]
Description=contestscore ham radio dashboard
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/contestscore
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
